# components/erlang/files/portmap.xml
<?xml version="1.0"?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">
<!--

 Copyright (c) 2020, Oracle and/or its affiliates. All rights reserved.

 CDDL HEADER START

 The contents of this file are subject to the terms of the
 Common Development and Distribution License (the "License").
 You may not use this file except in compliance with the License.

 You can obtain a copy of the license at usr/src/OPENSOLARIS.LICENSE
 or http://www.opensolaris.org/os/licensing.
 See the License for the specific language governing permissions
 and limitations under the License.

 When distributing Covered Code, include this CDDL HEADER in each
 file and include the License file at usr/src/OPENSOLARIS.LICENSE.
 If applicable, add the following below this CDDL HEADER, with the
 fields enclosed by brackets "[]" replaced with your own identifying
 information: Portions Copyright [yyyy] [name of copyright owner]

 CDDL HEADER END

 NOTE:  This service manifest is not editable; its contents will be
 overwritten by package or patch operations, including operating system
 upgrade.  Make customizations in a different file.

-->
<service_bundle type='manifest' name='portmap'>

  <service name='network/erlang/portmap' type='service' version='1' >

    <create_default_instance enabled='false' />
    <single_instance />

    <dependency name='fs' grouping='require_any' restart_on='none' type='service'>
      <service_fmri value='svc:/system/filesystem/minimal' />
    </dependency>

    <dependency name='network' grouping='require_any' restart_on='error' type='service'>
      <service_fmri value='svc:/milestone/network' />
    </dependency>

    <dependent name='rabbitmq' grouping='require_any' restart_on='none'>
      <service_fmri value='svc:/application/rabbitmq' />
    </dependent>

    <method_context>
      <method_credential user='_epmd' privileges='basic,{zone}:/system/volatile/_epmd' />
    </method_context>

    <exec_method type='method' name='start' exec='/lib/svc/method/portmap %m' timeout_seconds='5' />
    <exec_method type='method' name='stop' exec='/lib/svc/method/portmap %m' timeout_seconds='5' />

    <property_group name='config' type='application'>
      <propval name='port' type='integer' value='4369' />
      <property name='address' type='astring' />
      <propval name='logging' type='boolean' value='true' />
    </property_group>

    <stability value='Unstable' />

    <template>

      <common_name>
        <loctext xml:lang='C'>
        Erlang Port Mapper Daemon
        </loctext>
      </common_name>

      <documentation>
        <manpage title='epmd' section='1' />
      </documentation>

      <pg_pattern name='config' type='application' target='this'>
        <description><loctext xml:lang='C'>
            Erlang Port Mapper daemon configuration data as described in epmd(1).
         </loctext></description>
        <prop_pattern name='address' type='astring' required='false'>
          <description><loctext xml:lang='C'>
              List of comma-separated IP addresses on which EPMD daemon listens (default all).
          </loctext></description>
        </prop_pattern>
        <prop_pattern name='port' type='integer' required='false'>
          <description><loctext xml:lang='C'>
              TCP Port to which EPMD daemon listens (default 4369).
          </loctext></description>
          <cardinality min='0' max='1'/>
        </prop_pattern>
        <prop_pattern name='logging' type='boolean' required='false'>
          <description><loctext xml:lang='C'>
              Enable/disable logging to the system logs (default enabled).
          </loctext></description>
          <cardinality min='0' max='1'/>
        </prop_pattern>
      </pg_pattern>

    </template>

  </service>

</service_bundle>
